# fancontrol/fancontrol
!/bin/bash
# chkconfig: 2345 20 80
# description: Executes fancontrol. If it's properly configured, it'll
# control at least your cpu fan speed at startup

# Source function library.
. /etc/init.d/functions

case "$1" in 
    start)
       service fancontrol start
       ;;
    stop)
       service fancontrol stop
       ;;
    restart)
       service fancontrol stop
       service fancontrol start
       ;;
    status)
       # code to check status of app comes here 
       # example: status program_name
        service fancontrol status
       ;;
    *)                                                                                                                                                                                                     
       echo "Usage: $0 {start|stop|status|restart}"                                                                                                                                                        
esac                                                                                                                                                                                                       
                                                                                                                                                                                                           
exit 0  
